# launchd/com.clawdeck.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.clawdeck.server</string>
	<key>UserName</key>
	<string>claw</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/claw/.bun/bin/bun</string>
		<string>src/index.ts</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/claw/dev/clawdeck</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<!-- Évite une boucle de crash frénétique (ex. .env invalide au boot) :
	     launchd attend au moins 15 s entre deux relances. -->
	<key>ThrottleInterval</key>
	<integer>15</integer>
	<key>StandardOutPath</key>
	<string>/Users/claw/Library/Logs/clawdeck/stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/claw/Library/Logs/clawdeck/stderr.log</string>
</dict>
</plist>
